# armbian/userpatches/overlay/etc/systemd/system/atomicpi-hold-xmos.service
[Unit]
Description=Atomic Pi — release XMOS audio processor from reset
Documentation=https://www.digital-loggers.com/api_faqs.html
After=sysinit.target local-fs.target
Before=sound.target pulseaudio.service pipewire.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/atomicpi-xmos-reset
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
