# armbian/userpatches/overlay/etc/systemd/system/atomicpi-hold-mic.service
[Unit]
Description=Atomic Pi — hold XMOS microphone input selected (gpiochip1 line 0)
Documentation=https://www.digital-loggers.com/api_faqs.html
After=atomicpi-hold-xmos.service
Requires=atomicpi-hold-xmos.service

[Service]
Type=simple
ExecStart=/usr/bin/gpioset --hold gpiochip1/0=0
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
